# systemd/.config/systemd/user/sync_mail.service
[Unit]
Description=Sync mail.

[Service]
Type=oneshot
ExecStart=/home/chris/scripts/sync_mail.sh

[Install]
WantedBy=default.target
