# systemd/.config/systemd/user/sync_mail.timer
[Unit]
Description=Sync mail every 5 minutes.

[Timer]
OnStartupSec=0s
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
